# ops/systemd/frontend-ex.service
[Unit]
Description=frontend-ex (Phoenix SSR)
After=network.target

[Service]
Type=simple
# Adjust these three paths to match FX_DEPLOY_PATH (see deploy.sh).
WorkingDirectory=/opt/frontend-ex/current
EnvironmentFile=/opt/frontend-ex/frontend-ex.env
ExecStart=/opt/frontend-ex/current/bin/frontend_ex start
Restart=always
RestartSec=2
TimeoutStopSec=10
KillSignal=SIGTERM
SyslogIdentifier=frontend-ex
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
